# crontab
0 * * * * /home/sourabh/dev/spark/spark_etl.sh &> /home/sourabh/dev/spark/spark_etl.log

